# phase3/ollama.service.d/gotraceback.conf
# VM: /etc/systemd/system/ollama.service.d/gotraceback.conf
# Forces Go runtime to crash (and attempt core) on fault instead of only printing.
# Pair with coredump.conf (LimitCORE=infinity) and kernel core_pattern.
[Service]
Environment=GOTRACEBACK=crash
